# systemd/hf-gps-tec@.service
[Unit]
Description=HF PRN beacon recorder for ka9q-radio (%i)
Documentation=https://github.com/mijahauan/hf-gps-tec
After=network-online.target radiod@%i.service
Wants=network-online.target

[Service]
Type=notify
User=hfgpstec
Group=hfgpstec
ExecStart=/usr/local/bin/hf-gps-tec daemon \
    --config /etc/hf-gps-tec/hf-gps-tec-config.toml \
    --radiod-id %i \
    --instance %i
Restart=on-failure
RestartSec=15
TimeoutStopSec=30

# Logging — daemon process log goes to systemd journal; per-frame data
# goes to /var/lib/hf-gps-tec/<radiod_id>/.
StandardOutput=journal
StandardError=journal
SyslogIdentifier=hf-gps-tec@%i

# Environment — sigmond's coordination.env supplies per-radiod chain delays
# and runtime log-level overrides.
EnvironmentFile=-/etc/sigmond/coordination.env
Environment=PYTHONUNBUFFERED=1

# Sandboxing — standard recorder profile.
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
NoNewPrivileges=true
ReadWritePaths=/var/lib/hf-gps-tec /var/log/hf-gps-tec /var/lib/sigmond
RuntimeDirectory=hf-gps-tec
RuntimeDirectoryMode=0750

[Install]
WantedBy=multi-user.target
